# deployment/systemd/tsboard-api.service
[Unit]
Description=TSBOARD Go API Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=tsboard
Group=tsboard
WorkingDirectory=/var/www/tsboard
ExecStart=/var/www/tsboard/goapi-linux
Restart=always
RestartSec=5

# 환경변수 파일
EnvironmentFile=/var/www/tsboard/.env

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/tsboard/uploads /var/log/tsboard

# 리소스 제한
LimitNOFILE=65536
MemoryMax=512M

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsboard-api

[Install]
WantedBy=multi-user.target 